# systemd/fastblock.target
[Unit]
Description=fastblock target allowing to start/stop all fastblock*@.service instances at once

[Install]
WantedBy=multi-user.target